# guestsession/DEBIAN/postrm
#! /bin/sh

set -e
_DEB_SYSTEMD_HELPER_DEBUG=1

# In case this system is running systemd, we make systemd reload the unit files
# to pick up changes.
if [ -d /run/systemd/system ] ; then
	systemctl --system daemon-reload >/dev/null || true
fi
case "$1" in
    purge)
	#rm -f /etc/cups/printers.conf* /etc/cups/classes.conf*
	#rmdir /etc/cups 2>/dev/null || true
        if [ -x "/usr/bin/deb-systemd-helper-xe" ]; then
                #deb-systemd-helper-xe purge cups.path cups.service cups.socket >/dev/null
                deb-systemd-helper-xe purge guest-sesseion@Gast.service >/dev/null
                deb-systemd-helper-xe unmask guest-sesseion@Gast.service >/dev/null
        fi
	;;
    abort-upgrade)
#        if [ -e /etc/cups/cupsd.conf.conffile-bak ]; then
#            # Move cupsd.conf back
#            mv /etc/cups/cupsd.conf.conffile-bak /etc/cups/cupsd.conf
#        fi
        ;;
    remove)
        if [ -x "/usr/bin/deb-systemd-helper-xe" ]; then
                deb-systemd-helper-xe mask guest-sesseion@Gast.service >/dev/null
        fi
        ;;
    upgrade|failed-upgrade|abort-install|disappear)
        ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 0

esac
