# guestsession/DEBIAN/postinst
#!/bin/sh

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

# source debconf library
. /usr/share/debconf/confmodule

case "$1" in

  configure)
    #deb-systemd-helper-xe unmask guest-session@Gast.service >/dev/null || true
    systemctl unmask guest-session@Gast.service
    systemctl reenable --now guest-session@Gast.service
    #systemctl start guest-session@Gast.service

    ;;

  abort-upgrade|abort-remove|abort-deconfigure)
      exit 0
      ;;

  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
  ;;

esac

#echo "Noch was"
# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

##DEBHELPER#
#deb-systemd-helper-xe unmask guest-session@Gast.service >/dev/null || true
#
## was-enabled defaults to true, so new installations run enable.
#if deb-systemd-helper-xe --quiet was-enabled guest-session@Gast.service; then
#        # Enables the unit on first installation, creates new
#        # symlinks on upgrades if the unit file has changed.
#        deb-systemd-helper-xe enable guest-session@Gast.service >/dev/null || true
#else
#        # Update the statefile to add new symlinks (if any), which need to be
#        # cleaned up on purge. Also remove old symlinks.
#        deb-systemd-helper-xe update-state guest-session@Gast.service >/dev/null || true
#fi
#
#

#db_stop

exit 0
